# s6-overlay/s6-rc.d/apache2/run
#!/command/execlineb -P
/usr/sbin/apachectl -DFOREGROUND
